# client/eln-autofill.service
[Unit]
Description=ELN autofill run (labels, PubChem info, RDKit images)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=ELN_SERVER_URL=http://localhost:5000
Environment=ELN_API_KEY_FILE=/etc/eln-client/api_key
ExecStart=/usr/bin/python3 /opt/eln-client/eln_timer_client.py autofill
