# deploy/eln-server.service
[Unit]
Description=ELN server (Flask via gunicorn)
Wants=network-online.target
After=network-online.target

[Service]
# adjust User/paths to match the server
# server settings (Slack token path etc.) come from config.yaml in the checkout
WorkingDirectory=/usr/share/applications/eln_server
ExecStart=/usr/share/applications/eln_server/.venv/bin/gunicorn -w 1 -b 0.0.0.0:5000 app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
